# application/views/crontab
* * * * * /usr/bin/wget -O - -q -t 1 http://localhost:100/cron_job >> /var/www/html/CodeIgniter/application/controllers/logfile.txt
* * * * * /bin/date >> /tmp/cron_output

